# python-rtsp-service/rtsp-video.service
[Unit]
Description=RTSP Video Recording Service
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_GROUP
WorkingDirectory=/path/to/python-rtsp-service
ExecStart=/usr/bin/python3 /path/to/python-rtsp-service/video_service.py
Restart=always
RestartSec=10
StandardOutput=append:/path/to/python-rtsp-service/video_service.log
StandardError=append:/path/to/python-rtsp-service/video_service.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (optional)
# CPUQuota=200%
# MemoryLimit=2G

[Install]
WantedBy=multi-user.target
